# .udev/90-usb-keyboard.rules
ATTRS{idVendor}=="1532", ATTRS{idProduct}=="011c"
ACTION=="add", RUN+="/home/christian/.bin/keyboard_udev"
