# .systemd/avis_keyboard.service
[Unit]
Description=Swap keyboard layout

[Service]
ExecStart=/home/christian/.bin/swap_keyboard_mapping
